# .config/systemd/user/sync-math.service
[Unit]
Description="Timer job for syncing math files"

[Service]
ExecStart=/home/bigma/bin/sync-math
