# arch-setup/root/etc/systemd/system/sync-phone.service
[Unit]
Description="Timer job for syncing phone files"

[Service]
ExecStart=/usr/local/bin/sync-phone
